# crontab
*/30 * * * * feh --randomize --bg-fill /home/jens/pictures/wallpapers/
