# systemd/.config/systemd/user/dimmer.service
[Unit]
Description=Dimming based on time of day
Wants=dimmer.timer

[Service]
ExecStart=/usr/local/bin/dimmer