# systemd/.config/systemd/user/lowtmp.service
[Unit]
Description=Notify on low space in /tmp
Wants=lowtmp.timer

[Service]
ExecStart=/usr/local/bin/lowtmp