# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
Wants=mbsync.timer

[Service]
ExecStart=/usr/bin/mbsync -Va